# home/.config/systemd/user/fcitx5.service
[Unit]
Description=Fcitx5, a powerful Input Method Framework
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=fcitx5
