# home/.config/systemd/user/fcitx5.service
[Unit]
Description=Fcitx5, a powerful Input Method Framework
After=sway.service

[Service]
Type=simple
ExecStart=fcitx5
